# plugins/ship-flow/references/launchd/com.spacedock.ship-flow-scheduler.tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ship-flow-scheduler tick carrier (design.md §8). TEMPLATE — substitute
  @CONTROLLER_WORKTREE@ / @SPACEDOCK_BIN@ / @WORKFLOW_DIR@ before
  `launchctl load` (see docs/ship-flow/l3-scheduler-tick/RUNBOOK.md).
  The carrier owns scheduling only; all decisions live in the tick.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spacedock.ship-flow-scheduler.tick</string>
  <key>ProgramArguments</key>
  <array>
    <string>@CONTROLLER_WORKTREE@/plugins/ship-flow/bin/ship-flow-scheduler.sh</string>
    <string>tick</string>
    <string>--workflow-dir</string>
    <string>@WORKFLOW_DIR@</string>
    <string>--controller-worktree</string>
    <string>@CONTROLLER_WORKTREE@</string>
    <string>--runner</string>
    <string>gh</string>
    <string>--events-log</string>
    <string>@CONTROLLER_WORKTREE@/.ship-flow-scheduler-events.jsonl</string>
  </array>
  <key>WorkingDirectory</key>
  <string>@CONTROLLER_WORKTREE@</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>SPACEDOCK_BIN</key>
    <string>@SPACEDOCK_BIN@</string>
  </dict>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>@CONTROLLER_WORKTREE@/.ship-flow-scheduler-tick.log</string>
  <key>StandardErrorPath</key>
  <string>@CONTROLLER_WORKTREE@/.ship-flow-scheduler-tick.err.log</string>
</dict>
</plist>
